# configs/katashie-bot.service
[Unit]
Description=KATASHIE BOT — Plateforme SaaS de gestion de bots WhatsApp
Documentation=https://github.com/abesskamer237/katashie_bot.git
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/katashie-bot/backend
ExecStart=/usr/bin/node dist/index.js
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60

# Variables d'environnement
EnvironmentFile=/opt/katashie-bot/.env
Environment=NODE_ENV=production

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=katashie-bot

# Limites de ressources
LimitNOFILE=65535
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
